# deploy/systemd/brothershop.service
[Unit]
Description=BrotherShop CRM Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=brothershop
Group=brothershop
WorkingDirectory=/opt/brothershop
Environment=NODE_ENV=production
Environment=BOT_PYTHON_PATH=/opt/brothershop/.venv/bin/python
EnvironmentFile=/opt/brothershop/.env
ExecStart=/usr/bin/node /opt/brothershop/server.js
ExecStartPost=/bin/bash /opt/brothershop/scripts/systemd-healthcheck.sh
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

